# autoupgrade.timer
[Unit]
Description=Automatic Upgrade timer

[Timer]
OnCalendar=3:00
Persistent=false

[Install]
WantedBy=basic.target
