# autoupgrade.service
[Unit]
Description=AutoUpgrade

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /opt/upgrade.sh
